# infrastructure/php7-cli-supervisord/etc/supervisor/conf.d/messenger-consume.conf
[program:messenger-consunme]
command=php bin/console messenger:consume async --time-limit=3600
numprocs=5
process_name =%(program_name)s_%(process_num)02d
autostart=true
autorestart=true